# Meshwall/scripts/meshwall-listener.service
[Unit]
Description=MeshWall Active Listener Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/home/nyakki/Desktop/meshwall/venv/bin/python -m meshwall listen
Restart=on-failure
RestartSec=10
User=root
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/meshwall /var/log/meshwall /home/nyakki/.local/share/meshwall
Environment="MESHWALL_CONFIG=/etc/meshwall/meshwall.conf"

[Install]
WantedBy=multi-user.target