# Meshwall/scripts/meshwall-fetch.timer
[Unit]
Description=Daily MeshWall feed update
Requires=meshwall-fetch.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target